# toll-app.service
[Unit]
Description=E-ZPass Toll App (Gunicorn)
After=network.target

[Service]
Type=notify
User=YOUR_USER
Group=YOUR_GROUP
WorkingDirectory=/Users/ghuman/tolls
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/Users/ghuman/tolls"

# Gunicorn command
ExecStart=/Users/ghuman/tolls/venv/bin/gunicorn \
    --config /Users/ghuman/tolls/gunicorn_config.py \
    app:app

# Restart settings
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Security
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=toll-app

[Install]
WantedBy=multi-user.target

